# deliverance/paster-templates/deliverance/etc/supervisor.d/deliverance.conf
[program:deliverance]
command = ../bin/deliverance-proxy ./deliverance.xml
